# meta-mono/recipes-connectivity/uart-bridge/files/uart-bridge.service
[Unit]
Description=UART Bridge Daemon for STM32F411 Communication
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/uart-bridge
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
